# setupscripts/supervisor/gunicorn.conf
[program:gunicorn]
command=/etc/systemd/system/gunicorn.socket
directory=/home/djangoadmin/djexp/
user=djangoadmin
autostart=true
autorestart=true
redirect_stderr=true
